# systemd/user/pulseaudio.service
[Unit]
Description=PulseAudio Daemon
 
[Install]
WantedBy=default.target
 
[Service]
Type=simple
ExecStart=/usr/bin/pulseaudio

